# systemd/nats-server.internal.service
# nats-server.internal.service
[Unit]
Description=NATS Server (Internal Pub/Sub on port 4222)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nats
Group=nats
ExecStart=/usr/sbin/nats-server -p 4222
Restart=on-failure
RestartSec=2

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
