# systemd/nats-server.external.service
[Unit]
Description=NATS Server (External Pub/Sub on port 5222)
After=network.target
Wants=network.target

[Service]
Type=simple
User=nats
Group=nats
ExecStart=/usr/sbin/nats-server -p 5222
Restart=on-failure
RestartSec=2

# Security hardening (optional)
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
CapabilityBoundingSet=
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target


