# systemd/nats-server.external.service
# nats-server.external.service
[Unit]
Description=NATS Server (External Pub/Sub on port 5222)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nats
Group=nats
ExecStart=/usr/sbin/nats-server -p 5222
Restart=on-failure
RestartSec=2

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
